# container-files/etc/supervisord/sshd.conf
[program:sshd]
command=/usr/sbin/sshd -D
stdout_logfile=/var/log/supervisor/sshd.log, /dev/stdout
stderr_logfile=/var/log/supervisor/sshd.log, /dev/stdout
